# supervisor.conf
[program:novelio-queue-worker]
command=php /var/www/html/bulk-email-platform/artisan queue:work database --sleep=3 --tries=3 --timeout=120 --max-time=3600
directory=/var/www/html/bulk-email-platform
autostart=true
autorestart=true
startretries=5
numprocs=2
stopwaitsecs=120
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/novelio-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
